# backend/deploy/alps-backend.service
[Unit]
Description=ALPS Backend (read-only)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=alps
Group=alps
WorkingDirectory=/home/alps/alps/backend
EnvironmentFile=/home/alps/alps/backend/.env
ExecStart=/home/alps/.bun/bin/bun run src/index.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/alps /home/alps/alps/backend
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
